# instancesetup/crontab
*/5 * * * * /home/ubuntu/autoshutdown.sh
